# nagios/supervisord.conf
[supervisord]
nodaemon=true

[program:nagios]
command=/usr/sbin/nagios /etc/nagios/nagios.cfg

[program:apache2]
command=/usr/sbin/apachectl -D FOREGROUND
